# systemd/vinput-sherpa-zipformer.service
[Unit]
Description=Vinput Sherpa-ONNX WebSocket Server (Zipformer)
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/share/vinput/sherpa-onnx/bin/sherpa-onnx-offline-websocket-server \
  --port=18000 \
  --encoder=%h/.local/share/vinput/models/zipformer-zh-en/encoder-epoch-99-avg-1.onnx \
  --decoder=%h/.local/share/vinput/models/zipformer-zh-en/decoder-epoch-99-avg-1.onnx \
  --joiner=%h/.local/share/vinput/models/zipformer-zh-en/joiner-epoch-99-avg-1.onnx \
  --tokens=%h/.local/share/vinput/models/zipformer-zh-en/tokens.txt \
  --model-type=transducer \
  --provider=cpu \
  --num-threads=30 \
  --num-work-threads=5 \
  --max-batch-size=1 \
  --log-file=%h/.local/share/vinput/logs/zipformer-server.log
Restart=on-failure
RestartSec=5
Slice=user-expensive.slice

[Install]
WantedBy=default.target
