# systemd/vinput-sherpa-fire-red.service
[Unit]
Description=Vinput Sherpa-ONNX WebSocket Server (FireRed ASR)
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/share/vinput/sherpa-onnx/bin/sherpa-onnx-offline-websocket-server \
  --port=18001 \
  --fire-red-asr-encoder=%h/.local/share/vinput/models/sherpa-onnx-fire-red-asr2-zh_en-int8-2026-02-26/encoder.int8.onnx \
  --fire-red-asr-decoder=%h/.local/share/vinput/models/sherpa-onnx-fire-red-asr2-zh_en-int8-2026-02-26/decoder.int8.onnx \
  --tokens=%h/.local/share/vinput/models/sherpa-onnx-fire-red-asr2-zh_en-int8-2026-02-26/tokens.txt \
  --model-type=fire_red_asr \
  --num-threads=30 \
  --num-work-threads=5 \
  --max-batch-size=1 \
  --log-file=%h/.local/share/vinput/logs/fire-red-server.log
Restart=on-failure
RestartSec=5
Slice=user-expensive.slice

[Install]
WantedBy=default.target
